# Resources/SpookProvisioner/com.spookylabs.spooktacular.provisioner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    Spooktacular provisioner LaunchDaemon — the disk-injected
    variant.

    Not installed by any installer or pkg on the guest at all:
    `DiskInjector.installProvisionerDaemon` (host-side, root)
    attaches the VM's disk image with `diskutil image attach
    --nomount`, mounts the guest's data volume, and writes this
    plist straight into `/Library/LaunchDaemons/` (mode 0644,
    `root:wheel`) alongside the runner script in
    `/usr/local/libexec/` — all before the guest ever boots.
    `launchctl` picks the plist up on every subsequent boot;
    the runner looks for `first-boot.sh` on the per-VM virtio-fs
    share and executes it as root, then removes the trigger so
    subsequent boots no-op.

    ## Why a traditional absolute-path LaunchDaemon

    Both Spooktacular and its nested Guest Tools app are
    sandboxed. On macOS 14.4+, Apple forbids a sandboxed app
    from registering a LaunchDaemon via `SMAppService.daemon(
    plistName:)` unless the daemon itself is a sandboxed
    Mach-O — our runner is a bash script (needs to
    `mount_virtiofs` and exec arbitrary user scripts as root),
    which cannot meet that requirement. Writing the plist
    directly onto the guest's Data volume as a host-side root
    file operation, before first boot, sidesteps the sandboxed-
    daemon restriction entirely: there is no guest-side install
    step, no admin-password prompt, and nothing for the guest
    user to approve.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spookylabs.spooktacular.provisioner</string>

    <key>Program</key>
    <string>/usr/local/libexec/spook-provision-runner.sh</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
